# systemd/geoip.service
[Service]
ExecStart=%h/.local/bin/socket-geoip