# systemd/geoip.service
[Service]
ExecStart=%h/.cargo/bin/socket-geoip
Environment=TIMEOUT_SECS=60
Environment=DB_FILE=%h/GeoLite2-City.mmdb
